# extras/inference-server/examples/undertale-inference-worker.service
[Unit]
Description=Undertale Inference Worker
After=network.target

[Service]
ExecStart=inference worker --parallelism 4
Restart=on-failure

[Install]
WantedBy=multi-user.target
